# airootfs/etc/systemd/system/lightdm.service.d/override.conf
[Service]
# Ensure LightDM starts after gettingty on tty1
After=systemd-getty@1.service
# Don't start getty on graphical VTs
ExecStartPre=/bin/sh -c 'for vt in 2 3 4 5 6 7; do systemctl mask systemd-getty@$vt.service 2>/dev/null; done'
